# systemd/dci-openshift-agent.timer
[Unit]
Description=Trigger the DCI Openshift Agent

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
